# etc/systemd/system/timer-daily.target.wants/man-db.service
[Unit]
Description=Update manual page index caches

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/bin/mandb --quiet
