# etc/systemd/system/cups.service
.include /usr/lib/systemd/system/cups.service

[Unit]
PartOf=printer.target
