# etc/systemd/system/dimmer.service
[Unit]
Description=A simple backlight dimmer

[Service]
ExecStart=/usr/bin/dimmer --dim-by=100

[Install]
WantedBy=multi-user.target
