# etc/systemd/system/lxc-archlinux.service
[Unit]
Description=My little container

[Service]
ExecStart=/usr/bin/systemd-nspawn -jbD /srv/lxc-archlinux 3
KillMode=process
