# etc/systemd/system/lxc-archlinux.socket
[Unit]
Description=The SSH socket of my little container

[Socket]
ListenStream=23

[Install]
WantedBy=sockets.target
