# etc/systemd/system/mpd.socket
[Unit]
Description=mpd socket

[Socket]
ListenStream=6600
Service=mpd.service

[Install]
WantedBy=sockets.target
