# etc/systemd/system/net-auto-wireless.service
.include /usr/lib/systemd/system/net-auto-wireless.service

[Unit]
After=iptables.service
Requires=iptables.service

[Service]
KillMode=none
