# etc/systemd/system/net-auto-wired.service
.include /usr/lib/systemd/system/net-auto-wired.service

[Unit]
After=iptables.service
Requires=iptables.service

[Service]
TimeoutSec=2s
