# etc/systemd/system/netctl-ifplugd@.service
.include /usr/lib/systemd/system/netctl-ifplugd@.service

[Unit]
Before=network.target
Wants=network.target

[Service]
TimeoutSec=1s
