# etc/systemd/system/network.target
.include /usr/lib/systemd/system/network.target

[Unit]
Before=iptables.service
Requires=iptables.service
