# etc/systemd/system/slimlock.service
[Unit]
Description=Lock all ttys using vlock

[Service]
User=lahwaacz
ExecStart=/usr/bin/slimlock

[Install]
WantedBy=sleep.target
