# etc/systemd/system/sshd.service
.include /usr/lib/systemd/system/sshd.service

[Unit]
After=network.target
