# etc/systemd/system/timer-daily.timer
[Unit]
Description=Daily Timer

[Timer]
OnCalendar=*-*-* 21:30:00
Unit=timer-daily.target

[Install]
WantedBy=basic.target
