# etc/systemd/system/timer-weekly.timer
[Unit]
Description=Weekly Timer

[Timer]
OnCalendar=Sun *-*-* 21:30:00
Unit=timer-weekly.target

[Install]
WantedBy=basic.target
