# etc/systemd/system/timer-weekly.timer
[Unit]
Description=Weekly Timer

[Timer]
OnBootSec=30m
OnUnitActiveSec=7d
Unit=timer-weekly.target

[Install]
WantedBy=basic.target
